# meta-elisa-demo/recipes-elisa/vcan/files/can-dev-mapping-helper.service
[Unit]
Description=CAN interface helper

[Service]
ExecStart=/usr/sbin/can-dev-mapping-helper.sh
