# meta-elisa-demo/recipes-elisa/hello-world/files/helloworld.service
[Unit]
Description=Elisa Hello World

[Service]
Type=simple
ExecStart=/usr/bin/hello_world

[Install]
WantedBy=multi-user.target
