# meta-elisa-demo/recipes-elisa/pipe-demo/files/signal-source.service
[Unit]
Description=Elisa Demo Safety Application Signal Source

[Service]
Type=simple
ExecStart=/usr/bin/Signalsource

[Install]
WantedBy=multi-user.target

